# server/deploy/docferry-maintenance.service.example
[Unit]
Description=DocFerry backup, retention, and GC maintenance
After=network-online.target docker.service docferry-api.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/docferry/server
EnvironmentFile=/opt/docferry/server/.env.production
ExecStart=/opt/docferry/server/.venv/bin/python scripts/run_maintenance.py --apply --backup-root /root/backups/docferry --object-storage-root ${DOCFERRY_OBJECT_STORAGE_ROOT}
