# server/deploy/docferry-api.service.example
[Unit]
Description=Docferry FastAPI share server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/docferry/server
EnvironmentFile=/opt/docferry/server/.env.production
ExecStartPre=/opt/docferry/server/.venv/bin/alembic upgrade head
ExecStart=/opt/docferry/server/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8787
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
